# supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file = /var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:icinga2]
command = icinga2 daemon -x warning
autostart = true
autorestart = true

[program:icingaweb2]
command = httpd -DFOREGROUND
autostart = true
autorestart = true
